# roles/cloudflared_coolify/files/cloudflared-coolify.service
[Unit]
Description=Cloudflare Tunnel for coolify.designflow.app
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/cloudflared/coolify-tunnel.env
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate run --token ${TUNNEL_TOKEN}
Restart=always
RestartSec=5
TimeoutStopSec=15
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
