# roles/backrest_watchdog/files/backrest-dump-watchdog.timer
[Unit]
Description=Run the Backrest dump watchdog every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
