# roles/backrest_watchdog/templates/backrest-dump-watchdog.service.j2
# Managed by Ansible (backrest_watchdog role). Canonical source: backrest-wiz repo.
[Unit]
Description=Backrest backup-agent docker.sock staleness watchdog (plan §3.4)
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
Environment=BACKREST_CONTAINER={{ backrest_container_name }}
ExecStart=/usr/local/bin/backrest-dump-watchdog.sh
