# roles/backrest_watchdog/templates/backrest-dump-watchdog.timer.j2
# Managed by Ansible (backrest_watchdog role). Canonical source: backrest-wiz repo.
[Unit]
Description=Run the backrest docker.sock watchdog periodically (plan §3.4)

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ backrest_watchdog_interval }}
Persistent=true

[Install]
WantedBy=timers.target
