# roles/coolify/files/systemd/coolify-proxy-watchdog.service
[Unit]
Description=Coolify proxy docker-socket self-heal watchdog
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/coolify-proxy-socket-watchdog.sh
