# roles/coolify/files/systemd/coolify-autostart.service
[Unit]
Description=Auto-start Coolify containers stuck in Created state
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/local/bin/coolify-autostart.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
