# roles/dns_watchdog/files/dns-watchdog.service
[Unit]
Description=Repair DNS if external resolution is broken

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fix-dns-if-broken.sh
